# fs/rootfs/etc/inittab
# /etc/inittab
console::askfirst:-/bin/sh
::sysinit:/etc/init.d/rcS
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
::restart:/sbin/init
